# infrastructure/systemd/helios-ercot-congestion-batch.timer
[Unit]
Description=Run HeliosCTA ERCOT congestion support batch daily

[Timer]
OnCalendar=*-*-* 07:45:00 America/Chicago
Persistent=true
RandomizedDelaySec=10min
Unit=helios-ercot-congestion-batch.service

[Install]
WantedBy=timers.target
