# infrastructure/systemd/helios-weather-noaa-metar-observations.timer
[Unit]
Description=HeliosCTA NOAA AviationWeather METAR observation refresh timer

[Timer]
OnCalendar=*:07/15
Persistent=false
RandomizedDelaySec=2min
AccuracySec=1min
Unit=helios-weather-noaa-metar-observations.service

[Install]
WantedBy=timers.target
