# infrastructure/systemd/helios-isone-rt-hrl-scheduled-interchange.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Real-Time Hourly Scheduled Interchange orchestration timer

[Timer]
OnCalendar=*-*-* 06:25:00 UTC
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-rt-hrl-scheduled-interchange.service

[Install]
WantedBy=timers.target
