# infrastructure/systemd/helios-weather-wsi-hourly-forecast.service
[Unit]
Description=HeliosCTA WSI hourly forecast weather refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-weather-wsi-hourly-forecast.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.weather.wsi.hourly_forecast
TimeoutStartSec=45min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-weather-wsi-hourly-forecast

[Install]
WantedBy=multi-user.target
