# infrastructure/systemd/helios-pjm-data-miner-batch.service
[Unit]
Description=HeliosCTA PJM Data Miner scrape batch
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-pjm-data-miner-batch.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.pjm.data_miner_batch
TimeoutStartSec=6h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-pjm-data-miner-batch

[Install]
WantedBy=multi-user.target
