# infrastructure/systemd/helios-isone-hourly-system-demand.service
[Unit]
Description=HeliosCTA ISO-NE Hourly System Demand orchestration
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-isone-hourly-system-demand.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.isone.hourly_system_demand
TimeoutStartSec=2h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-isone-hourly-system-demand

[Install]
WantedBy=multi-user.target
