# infrastructure/systemd/helios-isone-da-hrl-cleared-demand.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Day-Ahead Hourly Cleared Demand orchestration timer

[Timer]
OnCalendar=*-*-* 17:20:00 UTC
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-da-hrl-cleared-demand.service

[Install]
WantedBy=timers.target
