# infrastructure/systemd/helios-isone-da-hrl-lmps.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Day-Ahead Hourly LMP orchestration timer

[Timer]
OnCalendar=*-*-* 13:10:00 America/New_York
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-da-hrl-lmps.service

[Install]
WantedBy=timers.target
