# infrastructure/systemd/helios-pjm-rt-unverified-hrl-lmps.service
[Unit]
Description=HeliosCTA PJM unverified hourly RT LMP refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-pjm-rt-unverified-hrl-lmps.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.pjm.rt_unverified_hrl_lmps
TimeoutStartSec=30min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-pjm-rt-unverified-hrl-lmps

[Install]
WantedBy=multi-user.target
