# infrastructure/systemd/helios-isone-hourly-system-demand.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Hourly System Demand orchestration timer

[Timer]
OnCalendar=*-*-* 06:10:00 UTC
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-hourly-system-demand.service

[Install]
WantedBy=timers.target
