# infrastructure/systemd/helios-isone-forecast-batch.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Forecast Batch orchestration timer

[Timer]
OnCalendar=*-*-* 15:20:00 UTC
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-forecast-batch.service

[Install]
WantedBy=timers.target
