# infrastructure/systemd/helios-ercot-load-batch.timer
[Unit]
Description=Daily HeliosCTA ERCOT load support scrape batch timer

[Timer]
OnCalendar=*-*-* 07:20:00 America/Chicago
Persistent=true
RandomizedDelaySec=10min
AccuracySec=1min
Unit=helios-ercot-load-batch.service

[Install]
WantedBy=timers.target
