# infrastructure/systemd/helios-ercot-renewables-5min-batch.timer
[Unit]
Description=Run HeliosCTA ERCOT 5-minute renewables actual support batch daily

[Timer]
OnCalendar=*-*-* 08:25:00 America/Chicago
Persistent=true
RandomizedDelaySec=10min
Unit=helios-ercot-renewables-5min-batch.service

[Install]
WantedBy=timers.target
