# infrastructure/systemd/helios-nav-positions.timer
[Unit]
Description=Daily HeliosCTA NAV positions scheduled SFTP poll timer

[Timer]
OnCalendar=*-*-* 04:00:00 America/Denver
Persistent=false
AccuracySec=1min
Unit=helios-nav-positions.service

[Install]
WantedBy=timers.target
