# infrastructure/systemd/helios-weather-wsi-daily-weighted-observations.service
[Unit]
Description=HeliosCTA WSI daily weighted observed weather refresh
Wants=network-online.target
After=network-online.target helios-weather-wsi-hourly-observed.service helios-weather-wsi-hourly-forecast.service helios-weather-wsi-daily-weighted-forecasts.service

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-weather-wsi-daily-weighted-observations.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.weather.wsi.daily_weighted_observations
TimeoutStartSec=45min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-weather-wsi-daily-weighted-observations

[Install]
WantedBy=multi-user.target
