# infrastructure/systemd/helios-pjm-hourly-price-backfill-7-day.timer
[Unit]
Description=Nightly HeliosCTA PJM price seven-day backfill repair timer

[Timer]
OnCalendar=*-*-* 02:00:00 America/New_York
Persistent=true
AccuracySec=5min
RandomizedDelaySec=10min
Unit=helios-pjm-hourly-price-backfill-7-day.service

[Install]
WantedBy=timers.target
