# infrastructure/systemd/helios-da-hrl-lmps.service
[Unit]
Description=HeliosCTA PJM Day-Ahead Hourly LMP orchestration
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-da-hrl-lmps.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.pjm.da_hrl_lmps
TimeoutStartSec=4h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-da-hrl-lmps

[Install]
WantedBy=multi-user.target
