# infrastructure/systemd/helios-ercot-dam-stlmnt-pnt-prices.service
[Unit]
Description=HeliosCTA ERCOT DAM Settlement Point Prices orchestration
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-ercot-dam-stlmnt-pnt-prices.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.ercot.dam_stlmnt_pnt_prices
TimeoutStartSec=2h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-ercot-dam-stlmnt-pnt-prices

[Install]
WantedBy=multi-user.target
