# infrastructure/systemd/helios-pjm-ops-sum.service
[Unit]
Description=HeliosCTA PJM Operations Summary scrape batch
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-pjm-ops-sum.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.pjm.ops_sum
TimeoutStartSec=2h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-pjm-ops-sum

[Install]
WantedBy=multi-user.target
