# infrastructure/systemd/helios-ercot-outage-capacity-batch.timer
[Unit]
Description=Run HeliosCTA ERCOT outage and capacity support batch daily

[Timer]
OnCalendar=*-*-* 13:35:00 UTC
Persistent=true
RandomizedDelaySec=10min
Unit=helios-ercot-outage-capacity-batch.service

[Install]
WantedBy=timers.target
