# infrastructure/systemd/helios-ercot-renewables-batch.timer
[Unit]
Description=Run HeliosCTA ERCOT renewables support batch daily

[Timer]
OnCalendar=*-*-* 08:10:00 America/Chicago
Persistent=true
RandomizedDelaySec=10min
Unit=helios-ercot-renewables-batch.service

[Install]
WantedBy=timers.target
