# infrastructure/systemd/helios-pjm-rt-hrl-lmps.timer
[Unit]
Description=Business-day HeliosCTA PJM verified hourly RT LMP polling timer

[Timer]
OnCalendar=Mon..Fri *-*-* 11:30:00 America/New_York
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-pjm-rt-hrl-lmps.service

[Install]
WantedBy=timers.target
