# infrastructure/systemd/helios-isone-forecast-batch.service
[Unit]
Description=HeliosCTA ISO-NE Forecast Batch orchestration
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-isone-forecast-batch.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.isone.forecast_batch
TimeoutStartSec=2h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-isone-forecast-batch

[Install]
WantedBy=multi-user.target
