# infrastructure/systemd/helios-ercot-congestion-batch.service
[Unit]
Description=HeliosCTA ERCOT congestion support batch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-ercot-congestion-batch.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.ercot.congestion_batch
Nice=5
NoNewPrivileges=true
