# infrastructure/systemd/helios-isone-rt-hrl-lmps-final.timer
[Unit]
Description=Daily HeliosCTA ISO-NE Final Real-Time Hourly LMP orchestration timer

[Timer]
OnCalendar=*-*-* 16:10:00 America/New_York
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=helios-isone-rt-hrl-lmps-final.service

[Install]
WantedBy=timers.target
