# infrastructure/systemd/helios-ercot-settlement-point-prices.service
[Unit]
Description=HeliosCTA ERCOT Real-Time Settlement Point Prices orchestration
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-ercot-settlement-point-prices.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.ercot.settlement_point_prices
TimeoutStartSec=45min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-ercot-settlement-point-prices

[Install]
WantedBy=multi-user.target
