# infrastructure/systemd/helios-prod-health-check.service
[Unit]
Description=HeliosCTA production health digest
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/opt/helioscta-platform/.venv/bin/python -m backend.orchestration.health.prod_health_check
TimeoutStartSec=5min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-prod-health-check

[Install]
WantedBy=multi-user.target
