# infrastructure/systemd/helios-pjm-hourly-bucket.timer
[Unit]
Description=Hourly HeliosCTA PJM scrape bucket timer

[Timer]
OnCalendar=*-*-* *:15:00 UTC
Persistent=false
RandomizedDelaySec=2min
AccuracySec=1min
Unit=helios-pjm-hourly-bucket.service

[Install]
WantedBy=timers.target
