# infrastructure/systemd/helios-slack-notification-outbox.service
[Unit]
Description=HeliosCTA Slack notification outbox sender
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-slack-notification-outbox.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.notifications.slack_outbox
TimeoutStartSec=5min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-slack-notification-outbox

[Install]
WantedBy=multi-user.target
