# infrastructure/systemd/helios-pjm-forecast-hourly.service
[Unit]
Description=HeliosCTA PJM hourly forecast refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=helios
Group=helios
WorkingDirectory=/opt/helioscta-platform
EnvironmentFile=/etc/helioscta/backend.env
ExecStart=/usr/bin/flock -n /tmp/helios-pjm-forecast-hourly.lock /opt/helioscta-platform/.venv/bin/python -m backend.orchestration.power.pjm.forecast_hourly
TimeoutStartSec=45min
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helios-pjm-forecast-hourly

[Install]
WantedBy=multi-user.target
